# contrib/systemd/electrumx-rpcd.service
[Unit]
Description=ElectrumX Rpc Http Server

[Service]
Type=simple
ExecStart=/path/to/electrumx/rpcd.py

[Install]
WantedBy=multi-user.target